# systemd/vackup_all.timer
[Unit]
Description=vackup_all service timer

[Timer]
OnCalendar=*-*-* 6:00:00
Unit=vackup_all.service

[Install]
WantedBy=timers.target